# besd.in
#!/bin/bash
#
# besd        Startup script for the OPeNDAP BES
#
# chkconfig: 2345 95 05
# description: OPeNDAP BES is a modular framework allowing access to
#	       data files
# processname: besctl
# config: @sysconfdir@/bes/bes.conf
# pidfile: @localstatedir@/run/bes.pid

# Path to the besctl script. Note: If the BES is installed using a prefix
# that is not set in the besctl script, add '-i <bes_prefix>' to the value
# of besctl here. This should only be needed in unusual situations.
prefix=@prefix@
exec_prefix=@exec_prefix@
besctl=@bindir@/besctl
prog=besd
RETVAL=0
LD_LIBRARY_PATH="@libdir@:${LD_LIBRARY_PATH}"
export LD_LIBRARY_PATH
# test -d "/var/run/bes" || mkdir -p /var/run/bes

# See how we were called.
case "$1" in
  start|stop|status|restart|help)
	what=${1}
	shift
	other=$@
	echo "${besctl} ${what} ${other}"
	${besctl} ${what} ${other}
	RETVAL=$?
	;;
  *)
	echo $"Usage: $prog {start|stop|restart|status|help}"
	exit 1
esac

exit $RETVAL
